# docker/php/crontabs/root
# Логи будут выводиться в docker (stdout)
* * * * * /usr/local/bin/php /var/www/html/artisan schedule:run

